# 常用脚本/定时数据库备份/auto_backup.cron
0 0 * * * sh /home/backup/auto_backup.sh>/dev/null 2>&1
